# deploy/showorsow-backend.service
# systemd unit for the Go backend. Install to /etc/systemd/system/.
# The binary is built by deploy.sh into /opt/showorsow/bin/backend.
[Unit]
Description=ShowOrSow Go backend (REST API + ledger writer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=showorsow
WorkingDirectory=/opt/showorsow
EnvironmentFile=/opt/showorsow/.env
ExecStart=/opt/showorsow/bin/backend
Restart=always
RestartSec=3
# Backend reads .env via EnvironmentFile above; pass an empty path so it does
# not also try to load a local .env file.
Environment=ENV_FILE=

[Install]
WantedBy=multi-user.target
