# deploy/showorsow-sandbox.service
# systemd unit for the self-contained Canton ledger (`dpm sandbox`).
# Install to /etc/systemd/system/. Runs as root because dpm lives under /root/.dpm.
# Persists to the sos-pg Docker Postgres (see canton-sandbox.conf). The Go backend
# + indexer connect to its JSON Ledger API v2 on 127.0.0.1:6864.
[Unit]
Description=ShowOrSow Canton sandbox (self-contained ledger, JSON API :6864)
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=root
Environment=HOME=/root
Environment=JAVA_OPTS=-Xmx4g
# Provides CANTON_PG_PASSWORD (Postgres store, see canton-sandbox.conf). The file
# is gitignored; other backend vars in it are simply ignored by Canton.
EnvironmentFile=/opt/showorsow/.env
WorkingDirectory=/opt/showorsow/canton
# NB: the `sandbox` command rejects --bootstrap, so DARs + the appOperator party
# are provisioned by sandbox-init.sh once the JSON API is up (idempotent).
ExecStart=/root/.dpm/bin/dpm sandbox \
  -c /opt/showorsow/deploy/canton-sandbox.conf \
  --no-tty \
  --log-file-name /opt/showorsow/canton/canton.log
ExecStartPost=/opt/showorsow/deploy/sandbox-init.sh
# Canton boots in ~40s and init waits for it; give it room.
TimeoutStartSec=240
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
